# midi-dashboard.conf
[program:app]
command=python3 app.py &
directory=/root/raspberry-pi-midi-dashboard
autostart=true
autorestart=true

[program:lcd]
command=python3 lcd.py &
directory=/root/raspberry-pi-midi-dashboard
autostart=true
autorestart=true

[group:midi_dashboard]
programs=app,lcd